# services/riak/config/etc/supervisor/conf.d/riak.conf
[program:riak]
command=/usr/sbin/riak console
pidfile=/var/log/riak/riak.pid
